# ops/ppgcomdata-backup.service
[Unit]
Description=PPGCOMDATA daily backup (Postgres + data + uploads)
Wants=network-online.target
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/ppgcomdata/backup.env
ExecStart=/usr/bin/env bash /root/projects/ppgcomdata/scripts/backup-vps.sh

# Defaults (override in /etc/ppgcomdata/backup.env)
Environment=BACKUP_ROOT=/root/backups/ppgcomdata
Environment=RETENTION_DAYS=14
Environment=GDRIVE_RETENTION_DAYS=30
Environment=PG_CONTAINER_NAME=ppgcomdata-db-prod
Environment=UPLOADS_VOLUME_NAME=prod_uploads

